# service.conf
[program:server]
command = /usr/local/bin/python -u /var/www/modules/Mailwarm/server.py
directory = /var/www
autostart = true
autorestart = true
startsecs = 0
startretries = 10000
restartpause = 5
